# data_sources_processing_src/crontab
0 1 * * * /usr/local/bin/python /app/update_datasets.py >> /data/dataprocessing.log 2>&1
